# contrib/aleesa-webapp.openrc
#!/sbin/openrc-run

name="$SVCNAME"
command="/var/lib/$SVCNAME/bin/$SVCNAME"
pidfile="/var/run/$SVCNAME/$SVCNAME.pid"

depend() {
	need net redis
	use dns
	after redis
}
